# roles/common/templates/etc/inittab.j2
# {{ ansible_banner }}

# What to do when Ctrl+Alt+Del is pressed.
ca:12345:ctrlaltdel:/bin/echo 'Ctrl+Alt+Del Disabled: /sbin/shutdown -t1 -a -r now'

